# systemd/xrdhttp@.socket
[Unit]
Description=XrdHttp socket

[Socket]
ListenStream=80
Service=xrootd@%i.service

[Install]
WantedBy=sockets.target
